# supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/dev/null
logfile_maxbytes=0
loglevel=info

[program:nginx]
command=nginx -g "daemon off;"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true
priority=1

[program:mailman]
command=/usr/local/bin/mailman
directory=/app
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true
priority=2
environment=SERVER_HOST=127.0.0.1,SERVER_PORT=8081,DB_DRIVER=sqlite,DB_NAME=/app/data/mailman.db,DB_HOST=localhost,DB_PORT=5432,DB_USER=mailman,DB_PASSWORD=,DB_SSLMODE=disable

[program:frontend]
command=node server.js
directory=/app/frontend
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true
priority=3
environment=PORT=3000,HOSTNAME=0.0.0.0